# part11/qa-project/qa-gateway.ini
[program:qa-gateway]
command=/data/www/gateway # 通过cargo build --release编译好的二进制文件gateway
directory=/data/www #二进制目录
autostart=true
startsecs=5
startretries=2
stdout_logfile=/data/logs/qa-gateway.log #
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
