# part11/qa-project/qa-svc.ini
[program:qa-svc]
command=/data/www/qa-svc # 通过cargo build --release编译好的二进制文件qa-svc
directory=/data/www # 二进制文件所在目录
autostart=true # 配置自动重启
startsecs=5 # 启动间隔时间，单位s
startretries=2 # 启动重试的次数
stdout_logfile=/data/logs/qa-svc.log # 日志目录
stdout_logfile_maxbytes=50MB # 日志文件大小限制
stdout_logfile_backups=10 # 日志文件备份个数
